# cron/crontab
*/20 * * * * /Users/mstockmatthews/code/rss/cron/run.sh
